# conf/crontab
* * * * * . $HOME/.profile; /usr/local/bin/php /var/www/html/src/cron.php >> /proc/1/fd/1 2>/proc/1/fd/2
# blah